# mobwrite/supervisord.conf
[program:mobwrite]
directory=/home/dotcloud/current/daemon
command=/home/dotcloud/current/daemon/mobwrite_daemon.py
